# config_file/frp.conf
[program:frp]
command =/home/nvidia/frp/frpc -c /home/nvidia/frp/frpc.ini
autostart=true
