# config_file/wificonnect.conf
[program:wificonnect]
command =nmcli connection up && nmcli connection modify FeiXingMoNin connection.autoconnect yes
autostart=true
